# dmibot.service
[Unit]
Description= Dmibot start service

[Service]
Type=simple
User=dmidev
ExecStart= /bin/bash /home/script/start.sh
